# infra/backup/sukhi-backup-verify.service
# SPDX-License-Identifier: AGPL-3.0-or-later
# Install to /etc/systemd/system/sukhi-backup-verify.service — see README.md.
[Unit]
Description=sukhi-fedi backup restore-verification (restic check + test-restore)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
EnvironmentFile=/etc/sukhi-fedi/backup.env
ExecStart=/usr/local/bin/sukhi-backup-verify.sh
